# scripts/backupscripts/crontar.cron
0 3 * * * /home/backup/backupscripts/tar-pkg.sh &>/dev/null
10 3 * * * /home/backup/backupscripts/tar-project.sh  &>/dev/null 
